# systemd/vllm-primary.service
[Unit]
Description=vLLM Primary — Llama 4 Scout 109B Q4_K_M (port 8080)
Documentation=https://github.com/mdl-systems/cocoro-llm-server
After=network.target nvidia-persistenced.service
Wants=nvidia-persistenced.service

# vllm-secondary は Primary に依存するため、Primary が先に起動
Before=vllm-secondary.service

[Service]
Type=simple
User=root
Group=root

# ---- 環境変数 ----
# .env はサービス起動スクリプト内で読み込むため、ここではシステム変数のみ定義
Environment="CUDA_HOME=/usr/local/cuda"
Environment="TORCH_CUDA_ARCH_LIST=12.0"
Environment="VLLM_ATTENTION_BACKEND=FLASHINFER"
Environment="VLLM_WORKER_MULTIPROC_METHOD=spawn"
Environment="CUDA_VISIBLE_DEVICES=0"
# Blackwell: NVLink非搭載のためP2PはデフォルトOFF（シングルGPU構成では無影響）
Environment="NCCL_P2P_DISABLE=0"

# PATH に CUDA と vLLM 仮想環境を含める
Environment="PATH=/opt/vllm-env/bin:/usr/local/cuda/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

# ---- 実行 ----
WorkingDirectory=/opt/cocoro-llm-server
ExecStart=/bin/bash /opt/cocoro-llm-server/vllm/start_primary.sh

# ---- ログ ----
# stdout/stderr は start_primary.sh 内で /var/log/cocoro-llm/vllm-primary.log に書き出し
# systemd journal にも並行して記録
StandardOutput=append:/var/log/cocoro-llm/vllm-primary.log
StandardError=append:/var/log/cocoro-llm/vllm-primary.log

# ---- リスタートポリシー ----
# OOM・モデルロードエラー含む異常終了時にリスタート
# ただし 3分以内に 5回失敗した場合は諦める（無限ループ防止）
Restart=on-failure
RestartSec=30
StartLimitIntervalSec=180
StartLimitBurst=5

# ---- タイムアウト ----
# Llama 4 Scout 109B のモデルロードは 96GB GPU でも最大 10分程度かかる
TimeoutStartSec=720
TimeoutStopSec=60

# ---- リソース制限 ----
# LimitNOFILE を引き上げ（vLLM は多くのファイルディスクリプタを使用）
LimitNOFILE=1048576
# GPU メモリ管理のため OOM スコアを下げる（システムより先に kill されにくくする）
OOMScoreAdjust=-500

# ---- セキュリティ ----
# GPU アクセスのため PrivateDevices は使えない
NoNewPrivileges=false
ProtectHome=false
ProtectSystem=false

[Install]
WantedBy=multi-user.target
